# src/php-ols/etc/s6-overlay/s6-rc.d/openlitespeed/run
#!/command/execlineb -P
with-contenv
/usr/local/lsws/bin/openlitespeed -d